# infra/systemd/crowdship-api.service
[Unit]
Description=Crowdship API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/crowdship
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /root/crowdship/src/server/http.js
Restart=on-failure
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
